# systemd/birdnet-api.service
[Unit]
Description=BirdNET-Pi read-only JSON API
After=network-online.target
Wants=network-online.target

[Service]
User=birdapi
Group=birdapi
Environment=BIRDS_DB=/home/birdnet/BirdNET-Pi/scripts/birds.db
ExecStart=/opt/birdnet-api/venv/bin/python /opt/birdnet-api/app/app.py
WorkingDirectory=/opt/birdnet-api/app
Restart=on-failure

[Install]
WantedBy=multi-user.target
