# systemd/birdviz.service
[Unit]
Description=BirdViz display
After=network-online.target
Wants=network-online.target

[Service]
Environment=BIRDVIZ_CONFIG=/opt/birdviz/configs/birdviz.yaml
WorkingDirectory=/opt/birdviz
ExecStart=/usr/bin/python3 -m birdviz.main
Restart=on-failure
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
